# src/data/org.buddiesofbudgie.BudgieDaemonV2.service.in
[Unit]
Description=Budgie Daemon v2
PartOf=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=org.buddiesofbudgie.BudgieDaemonV2
BusName=org.buddiesofbudgie.BudgieDaemon
Restart=on-failure
Slice=session.slice

[Install]
WantedBy=graphical-session.target